# 5.6/crontab
* * * * * echo "`date`: hello, world" >> /var/log/schedule.log
